# systemd/easymeter.timer
[Unit]
Description=easymeter auslesen alle 5min

[Timer]
OnActiveSec=30s
OnUnitActiveSec=2m

[Install]
WantedBy=basic.target
